# servers/supervisor/node.conf
[program:bloga]

command=sudo node /path/to/server/server.js
directory=/path/to/server
autostart=true
autorestart=true
stdout_logfile = /var/log/supervisor/app.log
stderr_logfile = /var/log/supervisor/app_err.log