# servers/upstart/app.conf
# appname

# description "A sample Upstart task for appname"
# author "Kevin Mann"

start on startup

script
    cd /home/kevin/ghost-0.3.3
    sudo npm start
end script